# deploy/etc/systemd/system/sshdgenkeys.service.d/epass-start-after-gui.conf
[Unit]
After=drm-arch-app.service

[Service]
ExecStartPre=/usr/local/bin/epass-wait-for-gui-ready.sh
