# deploy/etc/systemd/system/serial-getty@ttyS0.service.d/epass-fixed-local.conf
[Unit]
After=epass-boot-animation.service

[Service]
ExecStartPre=-/usr/local/bin/epass-boot-mark.sh serial-login:start
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- \\u' -L 115200 %I vt100
