# deploy/etc/systemd/system/drop-cache.timer
[Unit]
Description=Periodic cache drop for 64MB system

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
