# deploy/etc/systemd/system/drop-cache.service
[Unit]
Description=Drop kernel caches (64MB system)

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'sync && echo 3 > /proc/sys/vm/drop_caches'
