# deploy/etc/systemd/system/epass-gui-fallback.service
[Unit]
Description=ArkEPass GUI failure fallback
After=local-fs.target
RequiresMountsFor=/run /root

[Service]
Type=oneshot
ExecStart=/usr/local/bin/epass-gui-fallback.sh
StandardOutput=journal+console
StandardError=journal+console
